# share/systemd/omarchy-prayer-resume.service
[Unit]
Description=Rebuild omarchy-prayer timers after resume
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/omarchy-prayer-schedule

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
