# share/systemd/omarchy-prayer-schedule.service
[Unit]
Description=Omarchy prayer scheduler (rebuild today's timers)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/omarchy-prayer-schedule
